# examples/upstart/etc/init/sorry.conf
description "sorry server"
start on runlevel [2345]
stop on runlevel [!2345]
respawn
kill timeout 30
chdir /path/to/sorry
script
  . /etc/sysconfig/sorry
  exec /path/to/sorry/sorry
end script
